# etc/supervisord.d/confd_watch_nodes.ini
[program:confd_watch_nodes]
command=/root/bin/confd_watch_nodes.sh
redirect_stderr=true
stdout_syslog=true
