# etc/supervisord.d/slurmctld.ini
[program:slurmctld]
command=/root/bin/start_slurmctld.sh
redirect_stderr=true
stdout_syslog=true
